# deploy/systemd/shadowedvaca.service
[Unit]
Description=Shadowedvaca Site API
After=network.target postgresql.service

[Service]
User=www-data
WorkingDirectory=/opt/shadowedvaca
Environment="PYTHONPATH=/opt/shadowedvaca/src"
EnvironmentFile=/opt/shadowedvaca/.env
ExecStart=/opt/shadowedvaca/venv/bin/uvicorn sv_site.main:app \
    --host 127.0.0.1 \
    --port 8050 \
    --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
